# radio-monitor.service
[Unit]
Description=Radio to Spotify Monitor (all stations)
After=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/scripts/multi-radio
ExecStart=/home/pi/scripts/multi-radio/.venv/bin/python -W ignore::SyntaxWarning -m radio_monitor.main
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
